# deploy/license-server.service
[Unit]
Description=Angetic Essence License Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/angetic-license
EnvironmentFile=/opt/angetic-license/.env
ExecStart=/opt/angetic-license/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8080 \
    --timeout 120 \
    --access-logfile /var/log/angetic/access.log \
    --error-logfile /var/log/angetic/error.log \
    license_server.app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
